# allstar/debian/allstar-helpers.cron.d
# Get AllMondb on reboot
@reboot root /usr/sbin/astdb.php
