# allstar/debian/allstar-helpers.cron.d
# Get AllMon/Supermon/AllScan db on reboot
@reboot root /usr/sbin/astdb.php
