# allstar/debian/allstar-helpers.cron.d
# Get AllMondb on reboot
@reboot /usr/sbin/astdb.php
